# config/db/crontab
# Run backup daily at 2am
0 2 * * * /usr/local/bin/backup_postgres.sh >> /var/log/cron.log 2>&1
